# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/local/bin/mongod

[program:node] 
command=/bin/sh -c "cd /usr/local/start; node server"
